# tokens/com.openclaw.token-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.token-collector</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/matthew/.openclaw/workspace/tokens/run-collector.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/matthew/.openclaw/workspace/tokens</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/matthew/.openclaw/workspace/tokens/collector-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/matthew/.openclaw/workspace/tokens/collector-stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/matthew/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
